# systemd/99-boost-power.rules
# Boost Power Management - AC/Battery events
SUBSYSTEM=="power_supply", ATTR{type}=="Mains", ACTION=="change", RUN+="/usr/local/bin/ac-event"
